# launchd/com.mahos.codex-usage-tracker.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
  "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.mahos.codex-usage-tracker</string>
  <key>ProgramArguments</key>
  <array>
    <string>/opt/homebrew/bin/python3</string>
    <string>/Users/m/Documents/Codex/2026-05-19/i-d-like-a-little-codex/scripts/collect_codex_usage.py</string>
  </array>
  <key>RunAtLoad</key>
  <true/>
  <key>StartInterval</key>
  <integer>1800</integer>
  <key>WorkingDirectory</key>
  <string>/Users/m/Documents/Codex/2026-05-19/i-d-like-a-little-codex</string>
</dict>
</plist>
